# privaxy/debian/postrm
# !/bin/sh

#DEBHELPER#
if [ -x /etc/init.d/privaxy ]; then
    update-rc.d privaxy remove
    rm -f /etc/init.d/privaxy
fi

if systemctl is-enabled --quiet privaxy; then
    systemctl disable privaxy
fi

if systemctl status privaxy >/dev/null 2>&1 ; then
    systemctl stop privaxy
    systemctl disable privaxy
    rm -f /lib/systemd/system/privaxy.service
    systemctl daemon-reload
    systemctl reset-failed
fi

if id -u privaxy >/dev/null 2>&1 ; then
    deluser --quiet privaxy
    rm -rf /etc/privaxy
fi

exit 0